# forking-pinger.service
[Unit]
Description="Forking Pinger Daemon"

[Service]
Type=forking

ExecStart=/usr/local/bin/forking-pinger
