# oneshot-pinger.service
[Unit]
Description="Forking Ping Daemon"

[Service]
Type=oneshot

ExecStart=/home/root/forking-pinger
RemainAfterExit=yes
