# oneshot-pinger.service
[Unit]
Description="Forking Ping Daemon"

[Service]
Type=oneshot

ExecStart=/usr/local/bin/forking-pinger
RemainAfterExit=yes
